# util/linux_daemon.service
[Unit]
Description=Template daemon

[Service]
Type=forking
ExecStart=/usr/bin/linux_daemon

[Install]
WantedBy=multi-user.target

