# util/lxd.service
[Unit]
Description=Template daemon

[Service]
Type=forking
ExecStart=/usr/bin/lxd

[Install]
WantedBy=multi-user.target

